# scripts/pacpackages.service
[Unit]
Description=Updating the pacman index

[Service]
Type=simple
ExecStart=pacman -Sy
